# docker/router/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:inetsim]
directory=/etc/inetsim
command=/bin/bash /tmp/conf/run.sh
restart=always
